# supervisord.conf
[supervisord]
nodaemon=false

[program:tomcat]
priority=6
user=tomcat
command=catalina.sh run